# services/aerospike.service.conf
[program:aerospike]
user=root
command=/usr/local/bin/run-aerospike.sh
stdout_events_enabled = true
stderr_events_enabled = true
